# app-office/openerp/files/openerp
#!/sbin/runscript

command="/usr/bin/openerp-server"
command_args="-c ${CONFIGFILE}"
start_stop_daemon_args="--user ${USER} --group ${USER}"
background="true"
pidfile="${PIDFILE}"

depend() {
	after postgresql
}
