# app-office/openerp/files/openerp.initd
#!/sbin/runscript

pidfile="${PIDFILE}"
command="/usr/bin/openerp-server"
command_args="-c ${CONFIGFILE}"
command_background="true"
command_user="${USER}"
command_group="${USER}"

depend() {
	after postgresql
}

start_pre() {
	checkpath -q -d -m 0755 -o ${USER}:${GROUP} /var/run/openerp
}
