# app-admin/chef/files/initd/chef-client
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

CONFIGFILE=${CONFIGFILE:-/etc/chef/client.rb}
INTERVAL=${INTERVAL:-1800}
SPLAY=${SPLAY:-20}

pidfile=${CHEF_CLIENT_PIDFILE:-/var/run/chef/client.pid}
command=${CHEF_CLIENT_BINARY:-/usr/bin/chef-client}
command_args="-d -c ${CONFIGFILE} -i ${INTERVAL} -s ${SPLAY} -P ${pidfile} ${CHEF_CLIENT_OPTS}"

depend() {
	need net
}

start_pre() {
	checkpath -q -d -m 0755 -o chef:chef /var/run/chef
}
