# sys-process/dcron/files/dcron.service
[Unit]
Description=dcron

[Service]
ExecStart=/usr/sbin/crond -f

[Install]
WantedBy=multi-user.target
