# dev-db/couchdb/files/couchdb.init-0.10
#!/sbin/runscript
# Copyright 2008 Dirkjan Ochtman
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net epmd
}

checkconfig() {
	if test -z "$COUCHDB_PID_FILE"; then
		eerror "COUCHDB_PID_FILE must be set"
		return 1
	fi

	if test -z "$COUCHDB_USER"; then
		eerror "COUCHDB_USER must be set"
		return 1
	fi
}

start() {
	checkconfig || return 1

	args="-b -p $COUCHDB_PID_FILE "
	if test -n "$COUCHDB_STDOUT_FILE"; then
		args="$args -o $COUCHDB_STDOUT_FILE"
	fi
	if test -n "$COUCHDB_STDERR_FILE"; then
		args="$args -e $COUCHDB_STDERR_FILE"
	fi
	if test -n "$COUCHDB_RESPAWN_TIMEOUT"; then
		args="$args -r $COUCHDB_RESPAWN_TIMEOUT"
	fi
	if test -n "$COUCHDB_OPTIONS"; then
		args="$args $COUCHDB_OPTIONS"
	fi

	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start \
		--pidfile ${COUCHDB_PID_FILE} \
		--user ${COUCHDB_USER} \
		--exec ${EXEC} \
		-- $args
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop \
		--pidfile ${COUCHDB_PID_FILE}
	eend $?
}
