# net-analyzer/ganymed/files/ganymed.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Ganymed is an event collection daemon"

pidfile="/var/run/ganymed.pid"
config_file="/etc/ganymed/config.yml"
command="/usr/bin/ganymed"
command_args="-d -e production -p ${pidfile} -c ${config_file}"

depend() {
	need net
	use dns logger
}
